# src/logging/virtlogd-admin.socket.in
[Unit]
Description=Virtual machine log manager socket
Before=libvirtd.service

[Socket]
ListenStream=@localstatedir@/run/libvirt/virtlogd-admin-sock
Service=virtlogd.service

[Install]
WantedBy=sockets.target
