# src/logging/virtlogd-admin.socket.in
[Unit]
Description=Virtual machine log manager socket
Before=libvirtd.service
BindsTo=virtlogd.socket
After=virtlogd.socket

[Socket]
ListenStream=@runstatedir@/libvirt/virtlogd-admin-sock
Service=virtlogd.service
SocketMode=0600

[Install]
WantedBy=sockets.target
