# src/logging/virtlogd.socket.in
[Unit]
Description=libvirt logging daemon socket

[Socket]
ListenStream=@runstatedir@/libvirt/virtlogd-sock
Service=virtlogd.service
SocketMode=0600
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
Also=virtlogd-admin.socket
