# src/logging/virtlogd.socket.in
[Unit]
Description=Virtual machine log manager socket

[Socket]
ListenStream=@runstatedir@/libvirt/virtlogd-sock
Service=virtlogd.service
SocketMode=0600

[Install]
WantedBy=sockets.target
