# src/logging/virtlogd.socket.in
[Unit]
Description=Virtual machine log manager socket
Before=libvirtd.service

[Socket]
ListenStream=@localstatedir@/run/libvirt/virtlogd-sock

[Install]
WantedBy=sockets.target
