# src/nwfilter/virtnwfilterd.service.in
[Unit]
Description=Virtualization nwfilter daemon
Conflicts=libvirtd.service
Requires=virtnwfilterd.socket
Requires=virtnwfilterd-ro.socket
Requires=virtnwfilterd-admin.socket
After=network.target
After=dbus.service
After=apparmor.service
Documentation=man:virtnwfilterd(8)
Documentation=https://libvirt.org

[Service]
Type=notify
Environment=VIRTNWFILTERD_ARGS="--timeout 120"
EnvironmentFile=-@initconfdir@/virtnwfilterd
ExecStart=@sbindir@/virtnwfilterd $VIRTNWFILTERD_ARGS
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
Also=virtnwfilterd.socket
Also=virtnwfilterd-ro.socket
Also=virtnwfilterd-admin.socket
