# src/secret/virt-secret-init-encryption.service.in
[Unit]
Before=virtsecretd.service
Before=libvirtd.service
ConditionPathExists=!@localstatedir@/lib/libvirt/secrets/secrets-encryption-key

[Service]
Type=oneshot
UMask=0077
ExecStart=/usr/bin/sh -c 'dd if=/dev/random status=none bs=32 count=1 | systemd-creds encrypt --name=secrets-encryption-key - @localstatedir@/lib/libvirt/secrets/secrets-encryption-key'
