# src/node_device/virtnodedevd.service.in
[Unit]
Description=Virtualization nodedev daemon
Conflicts=libvirtd.service
Requires=virtnodedevd.socket
Requires=virtnodedevd-ro.socket
Requires=virtnodedevd-admin.socket
After=network.target
After=dbus.service
After=apparmor.service
Documentation=man:virtnodedevd(8)
Documentation=https://libvirt.org

[Service]
Type=notify
Environment=VIRTNODEDEVD_ARGS="--timeout 120"
EnvironmentFile=-@initconfdir@/virtnodedevd
ExecStart=@sbindir@/virtnodedevd $VIRTNODEDEVD_ARGS
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
Also=virtnodedevd.socket
Also=virtnodedevd-ro.socket
Also=virtnodedevd-admin.socket
