# daemon/virt-guest-shutdown.target
[Unit]
Description=Libvirt guests shutdown
Documentation=http://libvirt.org
