# src/vz/virtvzd.service.in
[Unit]
Description=Virtualization vz daemon
Conflicts=libvirtd.service
Requires=virtvzd.socket
Requires=virtvzd-ro.socket
Requires=virtvzd-admin.socket
After=network.target
After=dbus.service
After=apparmor.service
After=local-fs.target
After=remote-fs.target
Documentation=man:libvirtd(8)
Documentation=https://libvirt.org

[Service]
Type=notify
ExecStart=@sbindir@/virtvzd --timeout 120
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure

[Install]
WantedBy=multi-user.target
Also=virtvzd.socket
Also=virtvzd-ro.socket
Also=virtvzd-admin.socket
