# src/locking/virtlockd.socket.in
[Unit]
Description=Virtual machine lock manager socket

[Socket]
ListenStream=/var/run/libvirt/virtlockd/virtlockd.sock

[Install]
WantedBy=multi-user.target
