# src/locking/virtlockd.socket.in
[Unit]
Description=Virtual machine lock manager socket

[Socket]
ListenStream=@runstatedir@/libvirt/virtlockd-sock
Service=virtlockd.service
SocketMode=0600

[Install]
WantedBy=sockets.target
Also=virtlockd-admin.socket
