# src/locking/virtlockd-admin.socket.in
[Unit]
Description=libvirt locking daemon admin socket
BindsTo=virtlockd.socket
After=virtlockd.socket

[Socket]
ListenStream=@runstatedir@/libvirt/virtlockd-admin-sock
Service=virtlockd.service
SocketMode=0600
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
Also=virtlockd.socket
