# src/locking/virtlockd.socket.in
[Unit]
Description=libvirt locking daemon socket

[Socket]
ListenStream=@runstatedir@/libvirt/virtlockd-sock
Service=virtlockd.service
SocketMode=0600
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
Also=virtlockd-admin.socket
