# src/locking/virtlockd.service.in
[Unit]
Description=Virtual machine lock manager
Requires=virtlockd.socket
Before=libvirtd.service
Documentation=man:virtlockd(8)
Documentation=http://libvirt.org

[Service]
EnvironmentFile=-/etc/sysconfig/virtlockd
ExecStart=@sbindir@/virtlockd $VIRTLOCKD_ARGS
ExecReload=/bin/kill -USR1 $MAINPID
# Loosing the locks is a really bad thing that will
# cause the machine to be fenced (rebooted), so make
# sure we discourage OOM killer
OOMScoreAdjust=-900
# Needs to allow for max guests * average disks per guest
# libvirtd.service written to expect 4096 guests, so if we
# allow for 4 disks per guest, we get:
LimitNOFILE=16384

[Install]
Also=virtlockd.socket
