# src/locking/virtlockd.socket.in
[Unit]
Description=Virtual machine lock manager socket
Before=libvirtd.service

[Socket]
ListenStream=@runstatedir@/libvirt/virtlockd-sock
SocketMode=0600

[Install]
WantedBy=sockets.target
