# src/virtd-tls.socket.in
[Unit]
Description=libvirt @name@ daemon TLS IP socket
BindsTo=@service@.socket
After=@service@.socket
Conflicts=libvirtd-tls.socket
After=libvirtd-tls.socket

[Socket]
ListenStream=16514
Service=@service@.service

[Install]
WantedBy=sockets.target
