# src/remote/libvirtd.socket.in
[Unit]
Description=@name@ local socket
Before=@service@.service
@deps@

[Socket]
# The directory must match the @sysconfdir@/libvirt/@service@.conf unix_sock_dir setting
# when using systemd version < 227
ListenStream=@localstatedir@/run/libvirt/@sockprefix@-sock
Service=@service@.service
SocketMode=0666

[Install]
WantedBy=sockets.target
