# src/remote/libvirtd.socket.in
[Unit]
Description=@name@ local socket

[Socket]
ListenStream=@runstatedir@/libvirt/libvirt-sock
Service=libvirtd.service
SocketMode=@sockmode@
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
