# src/remote/libvirtd-ro.socket.in
[Unit]
Description=libvirt legacy monolithic daemon read-only socket
BindsTo=libvirtd.socket
After=libvirtd.socket

[Socket]
ListenStream=@runstatedir@/libvirt/libvirt-sock-ro
Service=libvirtd.service
SocketMode=0666
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
Also=libvirtd.socket
Also=libvirtd-admin.socket
