# src/remote/libvirtd-ro.socket.in
[Unit]
Description=@name@ local read-only socket
Before=@service@.service
BindsTo=@service@.socket
After=@service@.socket
@deps@

[Socket]
ListenStream=@runstatedir@/libvirt/@sockprefix@-sock-ro
Service=@service@.service
SocketMode=0666

[Install]
WantedBy=sockets.target
