# src/remote/libvirtd-admin.socket.in
[Unit]
Description=libvirt legacy monolithic daemon admin socket
BindsTo=libvirtd.socket
After=libvirtd.socket

[Socket]
ListenStream=@runstatedir@/libvirt/libvirt-admin-sock
Service=libvirtd.service
SocketMode=0600
RemoveOnStop=yes

[Install]
WantedBy=sockets.target
Also=libvirtd.socket
Also=libvirtd-ro.socket
