# src/remote/libvirtd-tls.socket.in
[Unit]
Description=@name@ TLS IP socket
Before=@service@.service
BindsTo=@service@.socket
After=@service@.socket
@deps@

[Socket]
# This must match the @sysconfdir@/libvirt/@service@.conf tls_port setting
# when using systemd version < 227
ListenStream=16514
Service=@service@.service

[Install]
WantedBy=sockets.target
