# src/remote/libvirtd-admin.socket.in
[Unit]
Description=@name@ admin socket
Before=libvirtd.service
BindsTo=libvirtd.socket
After=libvirtd.socket

[Socket]
ListenStream=@runstatedir@/libvirt/libvirt-admin-sock
Service=libvirtd.service
SocketMode=0600

[Install]
WantedBy=sockets.target
