# src/remote/libvirtd-admin.socket.in
[Unit]
Description=@name@ admin socket
Before=@service@.service
BindsTo=@service@.socket
After=@service@.socket
@deps@

[Socket]
# The directory must match the @sysconfdir@/libvirt/@service@.conf unix_sock_dir setting
# when using systemd version < 227
ListenStream=@localstatedir@/run/libvirt/@sockprefix@-admin-sock
Service=@service@.service
SocketMode=0600

[Install]
WantedBy=sockets.target
