# src/remote/virtproxyd.init.in
#!/sbin/openrc-run

description="Virtualization daemon proxy"

VIRTPROXYD_OPTS=${VIRTPROXYD_OPTS:-"${VIRTPROXYD_OPTS}"}
VIRTPROXYD_TIMEOUT=${VIRTPROXYD_TERMTIMEOUT:-"TERM/25/KILL/5"}

command="@sbindir@/virtproxyd"
command_args="-d ${VIRTPROXYD_OPTS}"
pidfile="@runstatedir@/virtproxyd.pid"
retry="${VIRTPROXYD_TERMTIMEOUT}"

extra_started_commands="reload"
description_reload="re-exec the daemon to enforce configuration reload"

depend() {
    use dbus
    after cgconfig corosync ebtables iptables ip6tables nfs nfsmount ntp-client ntpdportmap rpc.statd sanlock xenconsoled
    need net
}

reload() {
    ebegin "re-exec() virtproxyd"

    start-stop-daemon --signal SIGHUP \
      --exec "${command}" --pidfile "${pidfile}"
}
