# src/remote/virt-guest-shutdown.target
[Unit]
Description=Libvirt guests shutdown
Requires=libvirtd.service
Documentation=https://libvirt.org
