# src/remote/virt-guest-shutdown.target
[Unit]
Description=Libvirt guests shutdown
Documentation=https://libvirt.org
