# src/remote/virt-guest-shutdown.target
[Unit]
Description=libvirt guests shutdown target
Documentation=https://libvirt.org/
